# root/etc/services.d/garage_pie/run
#!/usr/bin/with-contenv sh

/app/garage_pie "${GARAGE_PIE_ARGS[@]}"
