# root/etc/services.d/garage_pie/run
#!/usr/bin/with-contenv sh

/app/normally-closed \
	--http-port 80 \
	"${NC_ARGS[@]}"
