# root/etc/services.d/normally_closed/run
#!/usr/bin/with-contenv sh

exec /app/normally-closed \
	--http-port 80 \
	/config/config.toml
